# docker/build/webserver/supervisord.conf
[supervisord]
nodaemon=true

[supervisorctl]


[program:apache2]
command=/opt/env-config/apache-start.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisord/apache.err.log
stdout_logfile=/var/log/supervisord/apache.out.log
